# alpine3.9.x-auto/_COPY_/etc/supervisor.d/05-bind.ini

[program:bind]
priority=5
command=/usr/sbin/named -u named -c /var/bind/bind.conf -f
directory=/var/bind
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_events_enabled=true
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stderr
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
numprocs=1
process_name=%(program_name)s

