# alpine3.9.x-auto/_COPY_/etc/supervisor.d/05-named.ini

[program:named]
priority=5
command=/usr/sbin/named -u named -c /var/bind/named.conf -f
process_name=%(program_name)s
numprocs=1
directory=/var/bind
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=false
stdout_logfile=/var/log/supervisord-named.log
stdout_logfile_maxbytes=1MB 
stdout_logfile_backups=10 
stdout_capture_maxbytes=1MB
