# alpine-openresty/_COPY_/etc/supervisor.d/20-redis.ini

[program:redis]
priority=20
command=/usr/bin/redis-server /etc/redis.conf
user=root
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_events_enabled=true
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
startsecs=3
