# alpine-openresty/_COPY_/etc/supervisor.d/40-nginx.ini

[program:nginx]
priority=40
command=nginx -g "daemon off;"
user=root
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_events_enabled=true
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
